# etc/periodic/hourly/000.zfs-snapshot
#!/bin/sh

# If there is a global system configuration file, suck it in.
#
if [ -r /etc/defaults/periodic.conf ]
then
    . /etc/defaults/periodic.conf
    source_periodic_confs
fi

pools=$hourly_zfs_snapshot_pools
if [ -z "$pools" ]; then
    pools='tank'
fi

keep=$hourly_zfs_snapshot_keep
if [ -z "$keep" ]; then
    keep=6
fi

case "$hourly_zfs_snapshot_enable" in
    [Yy][Ee][Ss])
        . /etc/periodic/zfs-snapshot
        do_snapshots "$pools" $keep 'hourly' "$hourly_zfs_snapshot_skip"
        ;;
    *)
        ;;
esac
